# .config/systemd/user/rsync-backup.service
[Unit]
Description=Rsync.net backup

[Service]
Type=simple

ExecStart=/home/phil/dotfiles/bin/backup.sh

